# logger-supervisor.ini
[program:logger-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/docker exec logger_app_1 php artisan queue:work redis --queue=default --tries=1 --daemon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=pavel
numprocs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/logger-worker.log
stopwaitsecs=3600

[program:logger-echo]
command=/usr/bin/docker exec logger_app_1 /usr/bin/laravel-echo-server start
autostart=true
autorestart=true
user=pavel
redirect_stderr=true
stdout_logfile=/var/log/supervisor/logger-echo.log
